# templates/postgresql-backup.timer.j2
[Unit]
Description=Run daily backup of Postgresql Server

[Timer]
{{ item.periodicity | default('{{ default_periodicity }}') }}
Persistent=true

[Install]
WantedBy=timers.target
