# templates/postgresql-backup.timer.j2
[Unit]
Description=Run daily backup of Postgresql Server

[Timer]
{{ item.periodicity | default(postgresql_backup_default_periodicity) }}
Persistent=true

[Install]
WantedBy=timers.target
